# mailhog/debian/prerm
#!/bin/bash

systemctl stop mailhog
systemctl disable mailhog
